# Buster/installer_images/build/payload/source/init_lcd.service
[Unit]
Description=Terastation Display Startup
Wants=network-online.target                                                                                  
After=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/ts3000_startup/startup_display

[Install]
WantedBy=multi-user.target
