# Stretch/installer_images/build/armhf-payload/source/micon_scripts/init_lcd.service
[Unit]
Description=Terastation Display Startup
Wants=network-online.target                                                                                  
After=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/micon_scripts/micon_startup_display

[Install]
WantedBy=multi-user.target
