# Stretch/installer_images/build/armhf-payload/source/micon_scripts/micon_boot.service
[Unit]
Description=Tell Terastation it finished booting so it doesnt shutdown
Wants=network-online.target                                                                                  
After=network-online.target

[Service]
ExecStart=/usr/local/bin/micon_scripts/micon_startup.py
Restart=on-failure
RestartSec=30

[Install]
WantedBy=multi-user.target
